# bundle/dot-claude/launchd/dev.ohmyclaude.resume-watchdog.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
oh-my-claude resume watchdog (macOS LaunchAgent).

Polls the resume-watchdog.sh script every 120 seconds to detect
unclaimed `resume_request.json` artifacts whose rate-limit windows
have cleared, and launches `claude --resume` in detached tmux
sessions.

Default: NOT loaded. Activate via `bundle/dot-claude/install-resume-watchdog.sh`
(or manually via `launchctl bootstrap gui/$UID ~/Library/LaunchAgents/dev.ohmyclaude.resume-watchdog.plist`).

The plist installer rewrites the placeholder __OMC_HOME__ and
__OMC_LOG_DIR__ tokens before loading.
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>dev.ohmyclaude.resume-watchdog</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>__OMC_HOME__/quality-pack/scripts/resume-watchdog.sh</string>
  </array>

  <!-- Poll cadence: every 120 seconds. The watchdog is stateless and
       cheap — a tick with nothing claimable exits in <100ms.

       KNOWN LIMITATION (v1.47 sre-lens): launchd's StartInterval is not
       robust across long laptop sleep / `bootout`. After a multi-day
       suspend, the interval timer can stop re-arming (observed in the
       field as `runs=1, pended nondemand spawn=interval`, heartbeat
       stale for 12 days). systemd's analogue uses `Persistent=true` to
       catch up on missed fires; launchd has no clean single-key
       equivalent (StartCalendarInterval would change the cadence
       semantics and needs a 60-element array for sub-hourly granularity;
       a conditional KeepAlive risks a relaunch storm because each tick
       exits immediately). The deliberate recovery is therefore NOT a
       plist change but the SessionStart self-heal in
       session-start-watchdog-health.sh: on the next Claude Code session
       (exactly when auto-resume matters — the user is back at the
       machine) the hook detects the stale heartbeat and re-registers the
       agent automatically (guarded against the orphan-spawn foot-gun).
       That recovery is platform-agnostic and covers ALL silent-death
       causes (sleep, bootout, macOS update, plist drift), so the plist
       is intentionally left as-is. -->
  <key>StartInterval</key>
  <integer>120</integer>

  <!-- Run at LaunchAgent load time so the user doesn't wait 2 minutes
       for the first tick after enabling. -->
  <key>RunAtLoad</key>
  <true/>

  <!-- Inherit the user's PATH so `claude` and `tmux` resolve. The
       installer baked in the user's actual $PATH at install time
       (preferring user-shell PATH when present, else the safe
       defaults below) — without this the daemon launches under
       launchd's barebones PATH and silently can't find claude when
       it's installed via npm / nvm / asdf / ~/.local/bin. -->
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>__OMC_PATH__</string>
    <key>HOME</key>
    <string>__OMC_USER_HOME__</string>
  </dict>

  <!-- Hard cap: even if the script hangs, kill it after 60s so the
       LaunchAgent slot is freed for the next tick. -->
  <key>ExitTimeOut</key>
  <integer>60</integer>

  <!-- Stdout / stderr land in a log file the user can tail to debug.
       The installer creates the log directory before loading. -->
  <key>StandardOutPath</key>
  <string>__OMC_LOG_DIR__/resume-watchdog.log</string>
  <key>StandardErrorPath</key>
  <string>__OMC_LOG_DIR__/resume-watchdog.err</string>

  <!-- Don't keep the daemon alive — each tick is a one-shot. -->
  <key>KeepAlive</key>
  <false/>

  <!-- Lower priority so the watchdog never preempts active work. -->
  <key>Nice</key>
  <integer>10</integer>
</dict>
</plist>
